# kitten-serve.service
[Unit]
Description=kitten-serve

[Service]
Type=simple
Restart=always
RestartSec=5s
ExecStart=/usr/local/bin/kitten-serve

[Install]
WantedBy=multi-user.target